# jankins_for_agetns/roles/jenkins_agent/templates/jenkins-agent1.service.j2
[Unit]
Description=Jenkins Agent 1
After=network.target

[Service]
Type=simple
User=jenkins
WorkingDirectory=/home/jenkins
ExecStart=/usr/bin/java -jar /home/jenkins/agent.jar -jnlpUrl http://your-jenkins-server/computer/agent1/slave-agent.jnlp -secret your-agent1-secret
Restart=always

[Install]
WantedBy=multi-user.target
